# supervisor/conf.d/websockets.conf
[program:websockets]
process_name=%(program_name)s
command=php /var/www/reactable/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
user=root